# .jobs/notes-sync.cron
*/5 * * * * ~/bin/auto-git-repo-sync ~/Documents/notes
